# systemd/atomic-queue.service
[Unit]
Description=Atomic Queue user service
After=default.target

[Service]
Type=simple
# Total in-memory queue cap, in integer GiB. Default 4. Override here
# or via a systemd drop-in (e.g. `systemctl edit --user atomic-queue`).
Environment=ATOMIC_QUEUE_MAX_QUEUED_GBYTES=4
ExecStart=%h/src/atomic-queue/atomic-queue serve --socket %t/atomic-queue/atomic-queue.sock
Restart=on-failure
RestartSec=1

[Install]
WantedBy=default.target
